# bin/sauron
#!/bin/sh
#
# chkconfig: - 96 16
# description: Aws cloudwatch daemon 
# processname: None
# config: None
# pidfile: None

# Source function library.
. /etc/init.d/functions

prog="sauron"
lockfile="/var/lock/subsys/${prog}"

# Exit if a single command fails
set -e

# Run as sauron
# Store the pid file in an appropriate place
PIDFILE="/var/run/${prog}.pid"
# Run sauron, currently with no options
DAEMON=sauron-daemon
DAEMON_OPTS=""

export PATH="${PATH:+$PATH:}/usr/sbin:/sbin"


function start() {
  # If the pid file exists
  if [ -e $PIDFILE ]; then
    echo "$prog is already running with pid `cat $PIDFILE`"
    exit 1
  fi
  # Start the thing, obviously
  echo -n $"Starting $prog: "
  nohup $DAEMON $DAEMON_OPTS > /dev/null 2>&1  &
  # wait a second
  sleep 1
  if [ -e $PIDFILE ]; then
    touch ${lockfile}
    echo_success
  else
    echo_failure
    exit 2
  fi
  echo
}

function stop() {
  # And stop it
  echo -n $"Stopping $prog: "
  if [ -e $PIDFILE ]; then
    kill `cat $PIDFILE` && echo_success || echo_failure
    rm -f ${lockfile}
    rm -f ${PIDFILE}
  else
    echo_sucess
  fi
  echo
}

case "$1" in
  start)
    # Invoke the start function
    start
    ;;
  stop)
    # Invoke the stop function
    stop
    ;;
  restart)
    # First invoke the stop function
    stop
    # And then start it
    start
    ;;
  *)
    # Didn't get any of the expected commands
    echo "Usage: "$1" {start|stop|restart}"
    exit 1
esac
