# files/sauron
#!/bin/sh
#
# chkconfig: - 96 16
# description: Aws cloudwatch daemon 
# processname: None
# config: None
# pidfile: None

# Source function library.
. /etc/init.d/functions

prog="sauron"
lockfile="/var/lock/subsys/${prog}"

# Run as sauron
# Store the pid file in an appropriate place
PIDFILE="/var/run/${prog}.pid"
# Run sauron, currently with no options
DAEMON=sauron-daemon
DAEMON_OPTS="--check ${prog}"
RETVAL=0

export PATH="${PATH}:/usr/sbin:/sbin:/usr/local/bin"


function start() {
  # If the pid file exists
  if [ -e $PIDFILE ]; then
    checkpid `cat $PIDFILE`
  else
    # Start the thing, obviously
    echo -n $"Starting $prog: "
    daemon $DAEMON_OPTS  "$DAEMON > /dev/null 2>&1 &"
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && touch /var/lock/subsys/$prog
  fi
  return $RETVAL
}

function stop() {
  # And stop it
  echo -n $"Stopping $prog: "
  killproc -p $PIDFILE $prog
  RETVAL=$?
  echo
  if [ $RETVAL -eq 0 ]; then
    rm -f /var/lock/subsys/${prog}
    rm -f /var/run/${prog}.pid
  fi
  return $RETVAL
}

function localstatus() {
  # show status
  echo -n $"show status of $prog: "
  status -p $PIDFILE $prog
  RETVAL=$?
  return $RETVAL
}

case "$1" in
  status)
    # Invoke the status function
    localstatus
    ;;
  start)
    # Invoke the start function
    start
    ;;
  stop)
    # Invoke the stop function
    stop
    ;;
  restart)
    # First invoke the stop function
    stop
    # And then start it
    start
    ;;
  *)
    # Didn't get any of the expected commands
    echo "Usage: "$1" {start|stop|status|restart}"
    exit 1
esac

exit $RETVAL
